# trunk/layout/DEBIAN/postinst
#!/bin/sh
/bin/launchctl load   /Library/LaunchDaemons/ch.ringwald.BTstack.plist 2&> /dev/null
exit 0
